# app/assets/cookbook/stunnel/client.service
[Unit]
Description={{cert.common_name}} SSL client

[Service]
ExecStart=/usr/bin/stunnel4 stunnel.conf
Restart=always
RestartSec=15s
WorkingDirectory=/opt/stunnel/{{cert.common_name}}/client


[Install]
WantedBy=multi-user.target